# cloudfront/update-ip-ranges.timer
[Unit]
Description=run update-ip-ranges daily

[Timer]
OnCalendar=*-*-* 00:00
RandomizedDelaySec=12h

[Install]
WantedBy=timers.target
